# ubuntu-setup.conf
[program:ubuntu]
command=/root/scripts/ubuntu.start.sh
numprocs=1
autostart=false
autorestart=false