# ubuntu-setup.conf
[program:ubuntu]
command=/root/scripts/ubuntu.start.sh
numprocs=1
autostart=false
autorestart=false

[program:cron]
command=cron -f
numprocs=1
autostart=true
autorestart=true